# templates/systemd/openclaw-vault-heartbeat-writer.service
[Unit]
Description=OpenClaw Vault Heartbeat Writer (manual oneshot)
Documentation=file:/opt/automation/projects/obsi-claw-AI_agent/docs/runbooks/OPENCLAW_VAULT_HEARTBEAT_WRITER_DEPLOY_SPRINT_4.md
After=local-fs.target
RequiresMountsFor=/opt/data/obsidian/vault-main /opt/automation/projects/obsi-claw-AI_agent
ConditionPathIsDirectory=/opt/data/obsidian/vault-main/Agent/Heartbeat
ConditionPathExists=/opt/automation/projects/obsi-claw-AI_agent/scripts/helpers/openclaw_vault_heartbeat_writer.py

[Service]
Type=oneshot
# La existencia real de este usuario y su membresia a `obsidian` queda pendiente de verificacion en host.
User=openclaw-vault-writer
Group=openclaw-vault-writer
SupplementaryGroups=obsidian
WorkingDirectory=/opt/automation/projects/obsi-claw-AI_agent
UMask=0007

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectClock=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any

LogsDirectory=openclaw-vault-heartbeat-writer
ReadWritePaths=/opt/data/obsidian/vault-main/Agent/Heartbeat
ReadOnlyPaths=/opt/automation/projects/obsi-claw-AI_agent/scripts/helpers/openclaw_vault_heartbeat_writer.py

ExecStart=/usr/bin/python3 /opt/automation/projects/obsi-claw-AI_agent/scripts/helpers/openclaw_vault_heartbeat_writer.py heartbeat.write --vault-root /opt/data/obsidian/vault-main --audit-root /var/log/openclaw-vault-heartbeat-writer --heartbeat-type runtime-status --run-id openclaw-vault-heartbeat-writer --source-ref docs/runbooks/OPENCLAW_VAULT_HEARTBEAT_WRITER_DEPLOY_SPRINT_4.md --context "systemd oneshot manual trigger" --result "runtime-status heartbeat created in Agent/Heartbeat"

StandardOutput=journal
StandardError=journal
